# 555-app.service
[Unit]
Description=555 Results Application
After=network.target mongodb.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/555-app
Environment="NODE_ENV=production"
Environment="PORT=3001"
ExecStart=/usr/bin/node /var/www/555-app/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
